# ws-signaling-server-and-test-client/openvpn-portmapio.service
[Unit]
Description=Manual OpenVPN connection using ash-temp-new.Server.ovpn.conf
After=network.target

[Service]
ExecStart=/usr/sbin/openvpn --config /etc/openvpn/client/ash-temp-new.Server.ovpn.conf
Restart=always
RestartSec=5
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
LimitNPROC=10
DeviceAllow=/dev/net/tun
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
